# deploy/systemd/symphony-trello@.service
[Unit]
Description=Symphony for Trello workflow %i
Documentation=https://github.com/martinfrancois/symphony-trello
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=symphony-trello
Group=symphony-trello
WorkingDirectory=/var/lib/symphony-trello
EnvironmentFile=-/etc/symphony-trello/service.env
LoadCredential=trello-api-key:/etc/symphony-trello/secrets/trello-api-key
LoadCredential=trello-api-token:/etc/symphony-trello/secrets/trello-api-token
ExecStart=/usr/bin/env java -jar /opt/symphony-trello/app/quarkus-run.jar /etc/symphony-trello/workflows/%i.WORKFLOW.md
Restart=on-failure
RestartSec=10
UMask=0077
CapabilityBoundingSet=
LockPersonality=true
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectProc=invisible
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ProtectHome=true
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
ReadWritePaths=/var/lib/symphony-trello
ReadOnlyPaths=/opt/symphony-trello /etc/symphony-trello

[Install]
WantedBy=multi-user.target
